# linux/h9-control.service
[Unit]
Description=Eventide H9 Control UI
After=network.target
Conflicts=zynthian.service zynthian-webconf.service jack2.service

[Service]
Type=simple
User=__INSTALL_USER__
WorkingDirectory=__PROJECT_DIR__
# Ensure we have a runtime directory for Wayland
Environment=XDG_RUNTIME_DIR=/run/user/__INSTALL_UID__
Environment="WLR_RENDERER_ALLOW_SOFTWARE=1"
# Start labwc (Wayland compositor) and have it run our startup script
ExecStart=/usr/bin/labwc -s __PROJECT_DIR__/linux/start-app.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
